# packaging/etherip-xdp@.service
# Templated systemd unit: one instance per external (uplink) interface.
#
#   sudo install -Dm0644 packaging/etherip-xdp@.service /etc/systemd/system/etherip-xdp@.service
#   sudo install -Dm0755 target/release/etherip-xdp /usr/bin/etherip-xdp
#   sudo mkdir -p /etc/etherip-xdp/interfaces.d/eth1   # one *.json per tunnel
#   sudo systemctl enable --now etherip-xdp@eth1
#   sudo systemctl reload etherip-xdp@eth1  # graceful config reload (SIGHUP)
#
# For the varlink management interface (etheripctl / etherip-xdp-manager), also
# enable the per-instance control socket and the host-wide manager:
#   sudo systemctl enable --now etherip-xdp-varlink@eth1.socket
#   sudo systemctl enable --now etherip-xdp-manager.socket
[Unit]
Description=XDP EtherIP tunnels on %i
Documentation=https://github.com/sorah/etherip-xdp
After=network-pre.target sys-subsystem-net-devices-%i.device
Wants=network-pre.target
BindsTo=sys-subsystem-net-devices-%i.device

[Service]
Type=simple
ExecStart=/usr/bin/etherip-xdp %i
# `systemctl reload` applies config changes gracefully (no data-path interruption
# for unaffected tunnels).
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2

# The varlink control socket fd is passed by name ("varlink"). The daemon also
# self-binds it when run outside socket activation, so this is optional.
Sockets=etherip-xdp-varlink@%i.socket

# Search $RUNTIME_DIRECTORY (/run/etherip-xdp) for volatile config drop-ins,
# ahead of /etc; preserved across restarts, cleared at reboot.
RuntimeDirectory=etherip-xdp
RuntimeDirectoryPreserve=yes

# Loading/attaching XDP, managing veth/neighbours, and ethtool need these.
# On kernels < 5.8, CAP_BPF does not exist; replace it with CAP_SYS_ADMIN.
# CAP_PERFMON sets allow_ptr_leaks so the verifier accepts the data-path's
# variable packet-offset arithmetic; without it newer kernels (6.x/7.0) reject
# the load with "pointer arithmetic ... prohibited for !root".
AmbientCapabilities=CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_NET_RAW

[Install]
WantedBy=multi-user.target
