# packaging/etherip-xdp@.service
# One instance per uplink interface:
#   sudo systemctl enable --now etherip-xdp@eth1
#   sudo systemctl reload etherip-xdp@eth1   # graceful config reload (SIGHUP)
# Enabling an instance also pulls in the per-uplink control socket and the
# host-wide manager socket, so etheripctl works without extra steps.
[Unit]
Description=XDP EtherIP tunnels on %i
Documentation=https://github.com/sorah/etherip-xdp
After=network-pre.target sys-subsystem-net-devices-%i.device
Wants=network-pre.target
BindsTo=sys-subsystem-net-devices-%i.device
# Pull in & order after the control socket; its fd is adopted via Sockets=.
Wants=etherip-xdp-varlink@%i.socket
After=etherip-xdp-varlink@%i.socket
# Pull in the host-wide manager socket so etheripctl works out of the box.
Wants=etherip-xdp-manager.socket

[Service]
Type=simple
ExecStart=/usr/bin/etherip-xdp %i
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2
Sockets=etherip-xdp-varlink@%i.socket

# Own dynamic user joining the static etherip-xdp-sock group, which owns the 0660
# sockets and /run/etherip-xdp so the manager (another dynamic user in the group)
# can reach them. A shared DynamicUser can't: systemd parks its RuntimeDirectory
# under /run/private 0700. RuntimeDirectory= is dropped for the same reason, so
# $RUNTIME_DIRECTORY (the config drop-in search path) is set explicitly.
DynamicUser=yes
SupplementaryGroups=etherip-xdp-sock
Environment=RUNTIME_DIRECTORY=/run/etherip-xdp

# CAP_SYS_ADMIN: unshare(CLONE_NEWNET) for the hidden veth peer (dropped by
# --disable-veth-peer-netns). CAP_PERFMON also sets allow_ptr_leaks for the
# data-path verifier. Pre-5.8: swap CAP_BPF for CAP_SYS_ADMIN.
AmbientCapabilities=CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_NET_RAW CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_NET_RAW CAP_SYS_ADMIN

# systemd-analyze security 2.2. AF_INET(6) is needed for the ethtool ioctl socket.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=net
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_PACKET AF_INET AF_INET6
SystemCallArchitectures=native
MemoryDenyWriteExecute=yes
IPAddressDeny=any
SystemCallFilter=@system-service bpf unshare setns
SystemCallFilter=~@keyring @resources @setuid @reboot @swap @cpu-emulation @obsolete
UMask=0077

[Install]
WantedBy=multi-user.target
