# packaging/etherip-xdp@.service
# Templated systemd unit: one instance per external (uplink) interface.
#
#   sudo install -Dm0644 packaging/etherip-xdp@.service /etc/systemd/system/etherip-xdp@.service
#   sudo install -Dm0755 target/release/etherip-xdp /usr/bin/etherip-xdp
#   sudo mkdir -p /etc/etherip-xdp/interfaces.d/eth1   # one *.json per tunnel
#   sudo systemctl enable --now etherip-xdp@eth1
#   sudo systemctl reload etherip-xdp@eth1  # graceful config reload (SIGHUP)
#
# For the varlink management interface (etheripctl / etherip-xdp-manager), also
# enable the per-instance control socket and the host-wide manager:
#   sudo systemctl enable --now etherip-xdp-varlink@eth1.socket
#   sudo systemctl enable --now etherip-xdp-manager.socket
[Unit]
Description=XDP EtherIP tunnels on %i
Documentation=https://github.com/sorah/etherip-xdp
After=network-pre.target sys-subsystem-net-devices-%i.device
Wants=network-pre.target
BindsTo=sys-subsystem-net-devices-%i.device

[Service]
Type=simple
ExecStart=/usr/bin/etherip-xdp %i
# `systemctl reload` applies config changes gracefully (no data-path interruption
# for unaffected tunnels).
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2

# The varlink control socket fd is passed by name ("varlink"). Under the
# sandboxing below the daemon runs unprivileged and cannot self-bind into the
# shared runtime dir (ProtectSystem=strict), so socket activation is the supported
# path for the management interface — enable etherip-xdp-varlink@%i.socket. The
# data path runs regardless of whether the socket is enabled.
Sockets=etherip-xdp-varlink@%i.socket

# Run as a dynamic, unprivileged user that joins the static etherip-xdp-sock group
# (created via sysusers.d). That group owns the 0660 varlink control sockets and
# the shared /run/etherip-xdp runtime dir (created via tmpfiles.d), so the
# host-wide manager — a separate dynamic user in the same group — can reach each
# daemon's per-device socket. (A shared DynamicUser identity can't: systemd parks
# its RuntimeDirectory under /run/private at 0700, unreachable cross-unit.)
DynamicUser=yes
SupplementaryGroups=etherip-xdp-sock

# Volatile config drop-ins are searched under $RUNTIME_DIRECTORY (/run/etherip-xdp),
# created by tmpfiles.d rather than RuntimeDirectory=: with DynamicUser= the latter
# would relocate it under /run/private (0700), out of the manager's reach and apart
# from the socket dir. Set the variable explicitly so the search still resolves.
Environment=RUNTIME_DIRECTORY=/run/etherip-xdp

# Capabilities (kept despite DynamicUser): CAP_BPF/CAP_PERFMON/CAP_NET_ADMIN/
# CAP_NET_RAW load & attach XDP, manage veth/neighbours and ethtool; CAP_PERFMON
# also sets allow_ptr_leaks so the verifier accepts the data-path's variable
# packet-offset arithmetic (6.x/7.0 reject it otherwise). CAP_SYS_RESOURCE: memlock.
# CAP_SYS_ADMIN: unshare(CLONE_NEWNET) for the hidden veth-peer netns (drop it with
# --disable-veth-peer-netns). On kernels < 5.8 replace CAP_BPF with CAP_SYS_ADMIN.
AmbientCapabilities=CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_NET_RAW CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_NET_RAW CAP_SYS_ADMIN

# Sandboxing (systemd-analyze security ~2.2 OK). The data path (XDP load, veth
# netns, ethtool, forwarding) is validated under this exact set:
#  - RestrictNamespaces=net keeps unshare(CLONE_NEWNET) working;
#  - RestrictAddressFamilies includes AF_INET(6) for the ethtool ioctl socket and
#    AF_NETLINK/AF_PACKET for netlink/XDP;
#  - SystemCallFilter keeps bpf/unshare/setns over @system-service.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=net
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_PACKET AF_INET AF_INET6
SystemCallArchitectures=native
MemoryDenyWriteExecute=yes
IPAddressDeny=any
SystemCallFilter=@system-service bpf unshare setns
SystemCallFilter=~@resources @reboot @swap @cpu-emulation @obsolete
UMask=0077

[Install]
WantedBy=multi-user.target
